# assets/com.ccx.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ccx.proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/backend-go/ccx</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__/backend-go</string>

    <!-- 手动启动模式：不随开机自启，需要时通过 ccx.sh start 启动 -->
    <key>RunAtLoad</key>
    <false/>

    <!-- 运行期间如果崩溃则自动重启，但不会开机自启 -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ccx.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ccx.stderr.log</string>
</dict>
</plist>
